# deploy/hermes-voice.service
[Unit]
Description=Hermes Voice Core
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dev/hermes-voice
EnvironmentFile=%h/.config/hermes-voice/environment
ExecStart=%h/.local/bin/uv run --frozen uvicorn hermes_voice.app:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
